# deployment/nutrition-bot.service
[Unit]
Description=FitAi Nutrition Tracker Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nutrition-bot
Environment="PATH=/opt/nutrition-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/opt/nutrition-bot/.env
ExecStart=/opt/nutrition-bot/venv/bin/python -m bot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
